# paas-worker/bocloud_worker/ansible_handler/playbooks/repository/telegraf/templates/systemd/system/telegraf.service.d/override.conf
[Service]
User={{ telegraf_runas_user }}
Group={{ telegraf_runas_group }}
